# client/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]

[program:sensu-client]
priority=10
directory=/tmp
command=/opt/sensu/bin/sensu-client -d /etc/sensu/conf.d -l /var/log/sensu/client.log
autostart=true
autorestart=true
